# crontab
30 7 * * 5 /usr/bin/python3 /app/main.py